# deploy/systemd/labweave-container-manager.service
[Unit]
Description=LabWeave Container Manager (toolchain resolver)
After=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/labweave/labweave.env
Environment=LISTEN_HOST=127.0.0.1
Environment=CORE_ROOT=%LABWEAVE_ROOT%
ExecStart=%LABWEAVE_ROOT%/.labweave/bin/container-manager
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
