# deploy/systemd/labweave-rule-engine.service
[Unit]
Description=LabWeave Rule Engine
After=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/labweave/labweave.env
Environment=LISTEN_HOST=127.0.0.1
WorkingDirectory=%LABWEAVE_ROOT%/rule-engine-py
ExecStart=%LABWEAVE_ROOT%/.venv/bin/python main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
